# sparse/lib/udev/rules.d/999-torch.rules
SUBSYSTEM=="leds", RUN+="/bin/chmod 666 /sys/class/leds/torch-light/brightness"
